# deploy/axfl-daily-runner.service
[Unit]
Description=AXFL Daily Runner (London+NY)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
# %h expands to the home directory of the specified User instance
WorkingDirectory=%h/Algotrading-Bot
EnvironmentFile=/etc/axfl/axfl.env
ExecStart=%h/Algotrading-Bot/scripts/run_daily_runner.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
